# examples/camera_viz/scripts/camera-streamer.service.in
# SPDX-FileCopyrightText: Copyright (c) 2026 NVIDIA CORPORATION & AFFILIATES. All rights reserved.
# SPDX-License-Identifier: Apache-2.0
#
# camera-streamer.service — installed by `camera_viz.sh deploy`.
# Placeholders are substituted at deploy time; do not hand-edit on the robot.

[Unit]
Description=Camera Streamer (camera_viz)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{WORKDIR}}
Environment=PATH={{VENV}}/bin:/usr/local/cuda/bin:/usr/bin:/bin
ExecStart={{VENV}}/bin/python {{WORKDIR}}/camera_streamer.py {{CONFIG}}
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
